# scripts/config_root/etc/systemd/system/onetime_init.service
[Unit]
Description=Generate SSH keys
Before=ssh.service

[Service]
Type=simple
ExecStart=/usr/bin/onetime_init

[Install]
WantedBy=multi-user.target
